# rob-the-bot.service
[Unit]
Description=Rob Discord Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=robbot
Group=robbot
WorkingDirectory=/opt/rob-the-bot/app
EnvironmentFile=/opt/rob-the-bot/app/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/rob-the-bot/app/.venv/bin/python /opt/rob-the-bot/app/main.py
Restart=always
RestartSec=10
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rob-the-bot/data
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
SystemCallArchitectures=native
SyslogIdentifier=rob-the-bot
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
